# conf/supervisor/production.ini
[program:gunicorn]
command=/home/ubuntu/environments/{{ project_name }}/bin/gunicorn wsgi:application -c /var/www/{{ project_name }}/conf/gunicorn/current.py
directory=/var/www/{{ project_name }}
user=nobody
group=nogroup
autostart=true
autorestart=true
redirect_stderr=True
stdout_logfile=/var/log/gunicorn.log

